# config/ctl/wordpress.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
redirect_stderr=true

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
redirect_stderr=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -F
autostart=true
numprocs=1
autorestart=true
redirect_stderr=true

[program:wordpress]
command=/usr/bin/php /path/to/your/wordpress/index.php --fpm
directory=/path/to/your/wordpress
user=yourusername
autostart=true
autorestart=true
redirect_stderr=true
